# nanoHUB/pipeline/dev.crontab
* * * * * make -C /app execute-notebook name=demo 2>&1 | logger -t cron
* * * * * sleep 10; make -C /app execute-notebook name=demo 2>&1 | logger -t cron
* * * * * sleep 20; docker run -a STDOUT --rm --name cron-test-3 pipeline_tasks python3 /app/hello-world.py
* * * * * sleep 30; make -C /app  execute-notebook name=demo 2>&1 | logger -t cron
* * * * * sleep 40; docker run -a STDOUT --rm --name cron-test-5 pipeline_tasks python3 /app/hello-world.py
* * * * * sleep 50; make -C /app  execute-notebook name=demo 2>&1 | logger -t cron
